# .devops/supervisord.dev.conf
; 开发容器 Supervisor 配置
; 包含: SSH + Xvfb + Chrome + CDP代理 + x11vnc + websockify

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=root:root

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
autoshutdown=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; SSH服务配置(开发容器核心,用于PyCharm远程调试)
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=5
startsecs=2

; Xvfb虚拟显示器配置
[program:xvfb]
command=bash -c "rm -f /tmp/.X1-lock && Xvfb :1 -screen 0 1280x1080x24"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10
startsecs=2

; Chrome浏览器配置
[program:chrome]
command=chromium \
    --display=:1 \
    --window-size=1280,1080 \
    --start-maximized \
    --no-sandbox \
    --disable-dev-shm-usage \
    --disable-setuid-sandbox \
    --disable-accelerated-2d-canvas \
    --disable-gpu \
    --disable-features=WelcomeExperience,SigninPromo \
    --no-first-run \
    --no-default-browser-check \
    --disable-infobars \
    --test-type \
    --disable-popup-blocking \
    --disable-gpu-sandbox \
    --no-xshm \
    --new-window=false \
    --disable-notifications \
    --disable-extensions \
    --disable-component-extensions-with-background-pages \
    --disable-prompt-on-repost \
    --disable-dialogs \
    --disable-modal-dialogs \
    --disable-web-security \
    --disable-site-isolation-trials \
    --remote-debugging-address=0.0.0.0 \
    --remote-debugging-port=8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DISPLAY=:1
priority=20
startretries=3
startsecs=5

; socat进程配置: 把Chrome内部CDP端口8222代理到外部9222
[program:socat]
command=socat TCP-LISTEN:9222,bind=0.0.0.0,fork,reuseaddr TCP:127.0.0.1:8222
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=2

; x11vnc服务配置(无密码,开发环境方便调试)
[program:x11vnc]
command=x11vnc -display :1 -nopw -shared -listen 0.0.0.0 -xkb -forever -rfbport 5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
startsecs=3

; Websockify配置: 将VNC转换为Websocket
[program:websockify]
command=websockify 0.0.0.0:5901 localhost:5900
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
startsecs=3

; 使用分组来统一启动or关闭进程
[group:services]
programs=sshd,xvfb,chrome,socat,x11vnc,websockify
